# supervisor.conf
[program:_TYPE_]
command=_SITE_ROOT_/start_server
directory=_SITE_ROOT_
redirect_stderr=true
stdout_logfile=_SITE_ROOT_/logs/_TYPE_.log
autorestart=true
